# launchd/com.user.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.automation</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{REPO_ROOT}}/scripts/run-automation.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>{{REPO_ROOT}}/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{REPO_ROOT}}/logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DOTNET_ROOT</key>
        <string>/opt/homebrew/opt/dotnet/libexec</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
